# systemd/gradi-vlcsync-gated.user.service
[Unit]
Description=Gradi vlcsync controller that auto-starts playback when all workers are ready
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HOSTS="192.168.0.4:5001 192.168.0.9:5001"
Environment=PATH=%h/.local/bin:/usr/bin:/bin
TimeoutStartSec=0
ExecStartPre=/bin/bash -lc 'for ep in $HOSTS; do ip=${ep%:*}; port=${ep#*:}; until nc -z -w1 "$ip" "$port"; do sleep 0.5; done; done'
ExecStart=%h/.local/bin/vlcsync --no-local-discovery \
  --rc-host 192.168.0.4:5001 \
  --rc-host 192.168.0.9:5001
ExecStartPost=/bin/bash -lc 'cmd=$'"'"'stop\nseek 0\nplay\n'"'"'; for ep in $HOSTS; do ip=${ep%:*}; port=${ep#*:}; printf "%b" "$cmd" | nc -N "$ip" "$port" & done; wait'
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
